# bash/cron_template.cron
* * * * * MASTER_KEY="p@r$3-mA&7er-k3y" APPLICATION_NAME="parseApplicationName" /root/update_song.sh